# satellite/deploy/nabu-satellite.service
# reSpeaker XVF3800 USB + HiFiBerry MiniAmp variant (the deployed shape) — provisioning rewrites
# the -D device on both audio commands below: the capture card is auto-detected by name from
# `arecord -l` (immune to ALSA index churn) and the speaker is the I2S DAC/amp HAT when one is
# present (`sndrpihifiberry*` card), else the mic card itself. The mic is 16 kHz mono native (no
# resampling; the XVF3800's 2 capture channels carry the same processed signal, so plughw's
# stereo->mono averaging is fine); plughw resamples only the 22050 Hz playback. The XVF3800's
# capture engine only runs while its OWN playback stream is active — provisioning detects this
# (capture-alone probe; sticky across re-provisions since the probe is only valid on a cold
# engine) and installs nabu-micclock.service, an endless zero-stream into its unconnected speaker
# out. The MiniAmp needs its overlay declared at provision time (DAC_OVERLAY=hifiberry-dac —
# one-shot config.txt write + reboot). No button, no LED by default.
# reSpeaker 2-Mic HAT variant: provision with the device arg plughw:CARD=seeed2micvoicec,DEV=0
# (16 kHz-native, used verbatim) and add --button-gpio 17 --led-spi (APA102s need dtparam=spi=on).
# Other hardware: a USB foot-switch is --button-evdev /dev/input/eventX:<keycode>; a wired
# indicator is --led-gpio <pin>.
# Music satellites (provision with MUSIC_HUB=<hub-host> MUSIC_ROOM=<room>) route ALL output —
# music (snapclient), TTS, and cues — to the speaker, mixed by PipeWire (neither the bcm2835 jack
# nor an I2S DAC can ALSA-dmix, so a userspace sound server does the mixing). Provisioning adds
# a systemd drop-in (nabu-satellite.service.d/pipewire.conf) that OVERRIDES the ExecStart below:
# --snd-command targets PipeWire through the `tts` ALSA softvol (aplay -D tts -> speaker), which
# holds the calibrated agent-voice level (provision-time TTS_VOLUME, default 65%) independent of
# music — the volume knob for amp HATs that have none (MiniAmp); --music-mixer Music --music-card
# <speaker> ducks the `music` softvol (snapclient) while active, and Environment=XDG_RUNTIME_DIR
# lets this system service reach the user's PipeWire (lingering is enabled).
# See scripts/provision-satellite-rs.sh. (Voice-only units use this ExecStart verbatim.)
[Unit]
Description=Nabu Rust Wyoming Satellite (%i)
After=network-online.target sound.target
Wants=network-online.target

[Service]
# Pin the performance governor: the steady sub-threshold wake-detection load sits exactly
# where ondemand idles at 600 MHz, roughly doubling per-chunk inference time and adding
# governor-ramp jitter to wake latency. A dedicated satellite Pi trades the extra idle power
# for consistent latency. The `+` prefix runs this as root despite User=%i; `|| true` keeps
# a governor-less kernel from failing the unit.
ExecStartPre=+/bin/sh -c 'for g in /sys/devices/system/cpu/cpufreq/policy*/scaling_governor; do echo performance > "$g" || true; done'
# Tail-latency insurance: background jobs (apt, unattended-upgrades) must not delay the
# 80 ms audio pumps.
Nice=-10
ExecStart=/usr/local/bin/nabu-satellite \
  --listen 0.0.0.0:10700 \
  --mic-command "arecord -D plughw:CARD=Array,DEV=0 -r 16000 -c 1 -f S16_LE -t raw -F 20000" \
  --snd-command "aplay -D plughw:CARD=sndrpihifiberry,DEV=0 -r 22050 -c 1 -f S16_LE -t raw --start-delay=100000 -F 50000" \
  --preroll-ms 1000 \
  --threshold 0.7 \
  --trigger-level 2
Environment=RUST_LOG=info
Restart=always
RestartSec=2
User=%i
# GPIO button/LED need the `gpio` group; a USB/evdev foot-switch needs `input`;
# the HAT's APA102 LEDs need `spi` (spidev nodes are root:spi on Raspberry Pi OS).
SupplementaryGroups=gpio input spi

[Install]
WantedBy=multi-user.target
